# bin/makerverse.service
[Unit]
Description=Makerverse Systemd service
After=docker.service

[Service]
Type=simple
WorkingDirectory=%h/makerverse
ExecStart=%h/makerverse/bin/launch
ExecReload=/bin/kill -HUP $MAINPID
Restart=on-failure
TimeoutStopSec=10

[Install]
WantedBy=multi-user.target
